# unifi-service/debian/unifi-core-agent-TEMPLATE/DEBIAN/postinst
#!/bin/bash

chown unifi: -R /opt/unifi

systemctl enable unifi-core-agent.service
systemctl start unifi-core-agent.service

