# unifi-service/debian/unifi-core-agent-TEMPLATE/DEBIAN/postinst
#!/bin/bash
set -euo pipefail

adduser --system unifiid

chown -R unifiid: /opt/unifi

systemctl enable unifi-core.service
systemctl start unifi-core.service
