# unifi-service/debian/unifi-core-agent-TEMPLATE/debian/postinst
#!/bin/bash

systemctl enable unifi-core-agent.service
systemctl start unifi-core-agent.service

