# unifi-service/debian/unifi-core-agent-TEMPLATE/debian/prerm
#!/bin/bash

systemctl stop unifi-core-agent.service
systemctl disable unifi-core-agent.service
systemctl reset-failed
