# unifi-service/debian/unifi-core-TEMPLATE/DEBIAN/postinst
#!/bin/bash

chown unifi: -R /opt/unifi

systemctl enable unifi-core.service
systemctl start unifi-core.service

