# unifi-service/debian/unifi-core-TEMPLATE/DEBIAN/postinst
#!/bin/bash

systemctl enable unifi-core.service
systemctl start unifi-core.service

