# unifi-service/debian/unifi-core-TEMPLATE/debian/prerm
#!/bin/bash

systemctl stop unifi-core.service
systemctl disable unifi-core.service
systemctl reset-failed
